# srcs/ftps/srcs/supervisord.conf
[supervisord]
nodaemon=true

[program:vsftpd]
command=sh run.sh
autostart=true
autorestart=true

[program:telegraf]
command=telegraf
autostart=true
autorestart=true
