# srcs/ftps/srcs/supervisord.conf
[supervisord]
nodaemon=true

[program:vsftpd]
command=vsftpd /etc/vsftpd/vsftpd.conf
autostart=true
autorestart=true

[program:telegraf]
command=telegraf
autostart=true
autorestart=true
