# srcs/wordpress/srcs/supervisord.conf
[supervisord]
nodaemon=true

[program:php-fpm7]
command=php-fpm7 --nodaemonize
autostart=true
autorestart=true

[program:nginx]
command=nginx -g 'daemon off;'
autostart=true
autorestart=true