# srcs/grafana/srcs/supervisord.conf
[supervisord]
nodaemon=true

[program:grafana]
command=./grafana-server
autostart=true
autorestart=true


[program:telegraf]
command=telegraf
autostart=true
autorestart=true
